# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:python]
command=/usr/bin/python /websocket_ssh.py

